# backend/koso.service
[Unit]
Description=Koso service
After=network.target

[Service]
Restart=always
RestartSec=1
TimeoutStopSec=30
User=root
WorkingDirectory=/root/koso/backend
Environment="KOSO_IMAGE_DIGEST=main"
ExecStart=docker run \
    --env KOSO_ENV=prod \
    -v /root/.secrets:/.secrets \
    --network=host \
    --rm \
    ghcr.io/kosolabs/koso@${KOSO_IMAGE_DIGEST}

[Install]
WantedBy=multi-user.target
